# config/supervisord.d/gunicorn.ini
[program:gunicorn]
command=gunicorn -c /root/code/myblog/config/gunicorn.conf wsgi:app
directory=/root/code/myblog/myblog
startsecs=5
stopwaitsecs=10
autostart=true
autorestart=true