# config/supervisord.d/gunicorn.ini
[program:gunicorn]
command=gunicorn -c /home/wenbao/myblog/config/gunicorn.conf wsgi:app
directory=/home/wenbao/myblog/myblog
startsecs=5
stopwaitsecs=10
autostart=true
autorestart=true