# config/supervisord.d/nginx.ini
[program:nginx]
command=/saas/softs/nginx/sbin/nginx
autostart=true
autorestart=true
startretries=3
exitcodes=0
stopsignal=KILL
stopwaitsecs=10
